# resin/data/bin/etc/update-motd.d/99-vpn-status
#!/bin/sh

if [ "$(wget -q -O - ipecho.net/plain)" = "173.29.176.159%" ]
then
	printf "THE VPN IS NOT ON!"

